# flarum/root/etc/cron.d/flarum
* * * * *  root  php /var/www/html/flarum schedule:run > /dev/null 2>&1